# xmadcoind.conf
description "xmadcoind daemon"
start on runlevel [23]
stop on shutdown
exec sudo  /home/azureuser/xmdcoind --no-console
post-stop exec sleep 30
respawn
respawn limit 5 30
